# Lab7.service
[Unit]
Description= Lab7

[Service]
Type=simple
